# files/daemon/git@.service
[Service]
User=gitolited
ExecStart=-/usr/libexec/git-core/git-daemon --interpolated-path=/var/lib/git/%H/%D --syslog --inetd
